# DFI2/hunter/dfi-dirty-capture.service
[Unit]
Description=DFI2 Dirty Traffic Capture — Watchlist/honeypot D2 from SPAN (24/7)
After=network.target clickhouse-server.service
Wants=clickhouse-server.service

[Service]
Type=simple
User=root
EnvironmentFile=/etc/dfi-hunter/env_dirty
ExecStart=/usr/bin/python3 -m hunter.hunter
WorkingDirectory=/opt/dfi2
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
Restart=always
RestartSec=5
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
